# conf/systemd/agoblinkm.service.in
[Unit]
Description=ago control blinkm i2c device

[Service]
ExecStart=@BINDIR@/agoblinkm
WorkingDirectory=@BINDIR@
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
