# conf/systemd/agodatalogger.service.in
[Unit]
Description=ago control Data Logger

[Service]
ExecStart=@BINDIR@/agodatalogger
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
