# conf/systemd/agoenigma2.service.in
[Unit]
Description=ago control Enigma2/Dreambox interface

[Service]
ExecStart=@BINDIR@/agoenigma2.py
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
