# conf/systemd/agojointspace.service.in
[Unit]
Description=ago control JointSpace Device

[Service]
ExecStart=@BINDIR@/agojointspace.py
WorkingDirectory=@BINDIR@
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
