# conf/systemd/agokwikwai.service.in
[Unit]
Description=ago control Kwikwai HDMI CEC Bridge Device

[Service]
ExecStart=@BINDIR@/agokwikwai
WorkingDirectory=@BINDIR@
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
