# conf/systemd/agomeloware.service.in
[Unit]
Description=ago control Meloware Lightswitchserver Gateway

[Service]
ExecStart=@BINDIR@/agomeloware.py
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
