# conf/systemd/agoowfs.service.in
[Unit]
Description=ago control 1-Wire Device

[Service]
ExecStart=@BINDIR@/agoowfs.py
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
