# conf/systemd/agorpc.service.in
[Unit]
Description=ago control HTTP RPC interface

[Service]
ExecStart=@BINDIR@/agorpc
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
