# conf/systemd/agoresolver.service.in
[Unit]
Description=ago control Resolver

[Service]
ExecStart=@BINDIR@/agoresolver
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
