# conf/systemd/agoscenario.service.in
[Unit]
Description=ago control Scenario Device

[Service]
ExecStart=@BINDIR@/agoscenario
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
