# conf/systemd/agotimer.service.in
[Unit]
Description=ago control Timer

[Service]
ExecStart=@BINDIR@/agotimer
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
