# conf/systemd/agox10.service.in
[Unit]
Description=ago control X10 Device

[Service]
ExecStart=@BINDIR@/agox10.py
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
